# systemd/myphotos-api.service.in
[Unit]
Description=MyPhotos API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@APP_USER@
Group=@APP_GROUP@
WorkingDirectory=@APP_DIR@
Environment=PYTHONUTF8=1
ExecStart=@APP_DIR@/.venv/bin/uvicorn app.api.main:app --host 127.0.0.1 --port 8888 --workers 2
Restart=on-failure
RestartSec=5
# Logs go to the systemd journal: journalctl -u myphotos-api -f

[Install]
WantedBy=multi-user.target
